# deploy/mac/templates/com.outerloop.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- HUB (coordinator) — binds loopback only; the ssh tunnel forwards the relay to it.
     Generated by the .pkg postinstall; do not edit in place (reinstall the pkg). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.outerloop.hub</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>inbox</string>
        <string>coordinator</string>
        <string>__HOST__</string>
        <string>8765</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__APPDIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>__PATH__</string>
        <key>HOME</key><string>__HOME__</string>
        <key>INBOX_HOME</key><string>__DATA__</string>
        <key>INBOX_FAKE</key><string>__FAKE__</string>
        <!-- hub's scheduler runs triage+scorer agents in real mode, so it needs claude too
             (launchd PATH won't find a per-user ~/.local/bin install). -->
        <key>INBOX_CLAUDE_BIN</key><string>__CLAUDE_BIN__</string>
        <!-- hub runs triage+scorer; override models here, e.g. "triage=haiku scorer=haiku" -->
        <key>INBOX_MODELS</key><string>__MODELS__</string>
    </dict>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>__DATA__/hub.log</string>
    <key>StandardErrorPath</key><string>__DATA__/hub.log</string>
    <key>ProcessType</key><string>Background</string>
</dict>
</plist>
