# deploy/mac/templates/com.outerloop.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- WORKER — runs on each Mac. Reaches the hub over HTTPS at the relay (or loopback
     on the mini itself). Generated by the .pkg postinstall. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.outerloop.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>inbox</string>
        <string>worker</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__APPDIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>__PATH__</string>
        <key>HOME</key><string>__HOME__</string>
        <key>INBOX_HOME</key><string>__DATA__</string>
        <key>INBOX_DEVICE</key><string>__DEVICE__</string>
        <key>INBOX_CAPABILITIES</key><string>__CAPS__</string>
        <key>INBOX_HUB</key><string>__HUB_URL__</string>
        <key>INBOX_DEVICE_TOKEN</key><string>__TOKEN__</string>
        <key>INBOX_FAKE</key><string>__FAKE__</string>
        <!-- Absolute claude path: launchd's PATH won't include a per-user ~/.local/bin
             install, so the postinstall resolves it and bakes it here (real mode). -->
        <key>INBOX_CLAUDE_BIN</key><string>__CLAUDE_BIN__</string>
        <!-- 1 = allow merging a repo with no CI configured (failing checks still block). -->
        <key>INBOX_ALLOW_MERGE_WITHOUT_CI</key><string>__ALLOW_NO_CI__</string>
        <!-- per-runner model overrides, e.g. "author=opus reviewer=opus"; empty = role defaults -->
        <key>INBOX_MODELS</key><string>__MODELS__</string>
    </dict>
    <!-- Restart on crash, but NOT after a clean exit: a worker with no hub configured
         exits 0 and stays stopped until the menu-bar Settings sets a hub and kickstarts it. -->
    <key>KeepAlive</key><dict><key>SuccessfulExit</key><false/></dict>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>__DATA__/worker.log</string>
    <key>StandardErrorPath</key><string>__DATA__/worker.log</string>
    <key>ProcessType</key><string>Background</string>
</dict>
</plist>
